# container-files/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/bin/nginx -g "daemon off;"
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.err.log
user=root